# 05-repoupdate.service
[Unit]
Description=Updates CentOS and EPEL yum repositories.

[Service]
Type=simple
ExecStart=/usr/local/sbin/repoupdate
StandardOutput=syslog